# usr/lib/systemd/system/post-temperature.service
[Unit]
Description=Post temperature to student server

[Service]
Type=oneshot
ExecStart=/usr/bin/post-temperature
